# archiso/airootfs/etc/skel/.config/autostart/frog-init.desktop
[Desktop Entry]
Type=Application
Name=Frog Init
Exec=/usr/local/bin/frog-init.sh
X-GNOME-Autostart-enabled=true
NoDisplay=true
